# nextcloud/rootfs/etc/s6-overlay/s6-rc.d/init-postgres-initdb/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
set -e

# Read db_password from HA options if available
OPTIONS_FILE="/data/options.json"
if [[ -f "${OPTIONS_FILE}" ]]; then
    OPT_DB_PASS=$(jq -r '.db_password // empty' "${OPTIONS_FILE}")
    if [[ -n "${OPT_DB_PASS}" ]]; then
        export POSTGRES_PASSWORD="${OPT_DB_PASS}"
    fi
fi

POSTGRES_USER="${POSTGRES_USER:-nextcloud}"
POSTGRES_PASSWORD="${POSTGRES_PASSWORD:-nextcloud}"
POSTGRES_DB="${POSTGRES_DB:-nextcloud}"

# Check if database is already initialized
if [ -s "${PGDATA}/PG_VERSION" ]; then
    echo "PostgreSQL database already initialized in ${PGDATA}"
    # Ensure ownership is correct
    chown -R postgres:postgres "${PGDATA}"
    exit 0
fi

echo "Initializing PostgreSQL database..."

# Initialize the database cluster
s6-setuidgid postgres initdb \
    --username="${POSTGRES_USER}" \
    --pwfile=<(printf "%s\n" "${POSTGRES_PASSWORD}") \
    --encoding=UTF8 \
    --locale=C \
    -D "${PGDATA}"

# Configure pg_hba.conf for local connections
cat >> "${PGDATA}/pg_hba.conf" <<-EOF

# Allow local connections for nextcloud
host all all 127.0.0.1/32 md5
host all all ::1/128 md5
local all all md5
EOF

# Configure postgresql.conf for local-only access
sed -ri "s!^#?(listen_addresses)\s*=\s*\S+.*!\1 = 'localhost'!" "${PGDATA}/postgresql.conf"

# Start temporary server to create the database
s6-setuidgid postgres pg_ctl -D "${PGDATA}" -o "-c listen_addresses='localhost'" -w start

# Create the nextcloud database if it doesn't exist
s6-setuidgid postgres psql -v ON_ERROR_STOP=1 --username="${POSTGRES_USER}" <<-EOSQL
    SELECT 'CREATE DATABASE ${POSTGRES_DB}' WHERE NOT EXISTS (SELECT FROM pg_database WHERE datname = '${POSTGRES_DB}')\gexec
EOSQL

# Stop temporary server
s6-setuidgid postgres pg_ctl -D "${PGDATA}" -m fast -w stop

echo "PostgreSQL database initialization complete"
